# supervisor.conf
[program:jhongbot]
command=wrapper.sh start
directory=%(here)s
autostart=true
autorestart=true
startretries=3
stderr_logfile=%(here)s/error.log
stdout_logfile=%(here)s/jhongbot.log
user=ubuntu